# touchscreen-detector.service
[Unit]
Description=Touchscreen Detection Service
After=graphical.target
Wants=graphical.target
Requires=graphical.target

[Service]
Type=simple
User=%i
Group=%i
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
WorkingDirectory=/home/%i/touchscreen-detector
ExecStart=/usr/bin/python3 /home/%i/touchscreen-detector/detecttouch.py
Restart=no
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=touchscreen-detector

# Ensure the service has access to input devices
SupplementaryGroups=input

[Install]
WantedBy=graphical.target 