# touchscreen-detector.service
[Unit]
Description=Touchscreen Detection Service
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
User=root
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/%i/.Xauthority
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=XDG_SESSION_TYPE=x11
Environment=XDG_CURRENT_DESKTOP=XFCE
Environment=DESKTOP_SESSION=xubuntu
WorkingDirectory=/home/%i/touchscreen-detector
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/python3 /home/%i/touchscreen-detector/detecttouch.py
Restart=no
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=touchscreen-detector
KillMode=mixed

# Ensure the service has access to input devices and audio
SupplementaryGroups=input audio video

[Install]
WantedBy=graphical.target 